# etc/rc.d/S99hwf_info_upload
#!/bin/sh /etc/rc.common

START=99

start() {
	hwf_info_upload.sh &
}

stop() {
	killall hwf_info_upload.sh &> /dev/null
}
