# usr/lib/opkg/info/autossh.list
/etc/config/autossh
/etc/init.d/autossh
/usr/sbin/autossh
