# usr/lib/opkg/info/hwfatd.list
/usr/sbin/hwfatd
/etc/init.d/hwfatd
/usr/sbin/hwf-at
