# usr/lib/opkg/info/mentohust.list
/etc/config/mentohust
/etc/init.d/mentohust
/usr/sbin/mentohust
