# usr/lib/opkg/info/syslog-ng3.list
/etc/hotplug.d/block/restart.d/80-syslog-ng
/etc/init.d/syslog-ng
/etc/syslog-ng.conf_sample
/usr/sbin/syslog-ng
